# distrib/coreutils/rc.local
#!/bin/sh

echo 
echo Welcome to Xv6!
echo 
